# files/etc/uci-defaults/93-autoreboot-enable
#!/bin/sh

[ -x /etc/init.d/autoreboot ] || exit 0

/etc/init.d/autoreboot enable >/dev/null 2>&1 || true
/etc/init.d/autoreboot start >/dev/null 2>&1 || true

exit 0
