# files/etc/uci-defaults/94-tailscale-settings-enable
#!/bin/sh

[ -x /etc/init.d/tailscale-settings ] || exit 0

/etc/init.d/tailscale-settings enable >/dev/null 2>&1 || true
/etc/init.d/tailscale-settings start >/dev/null 2>&1 || true

exit 0
