# files/etc/uci-defaults/95-tailscale-accept-dns-guard
#!/bin/sh

/etc/init.d/tailscale-accept-dns-guard enable >/dev/null 2>&1 || true
/etc/init.d/tailscale-accept-dns-guard start >/dev/null 2>&1 || true

exit 0
