# files/etc/uci-defaults/99-tailscale-quad100-health
#!/bin/sh

/etc/init.d/tailscale-quad100-health enable >/dev/null 2>&1 || true
/etc/init.d/tailscale-quad100-health start >/dev/null 2>&1 || true

exit 0
