# package/luci-app-podman/root/etc/cron.d/podman-cleanup
# Cleanup orphaned Podman log stream sessions (older than 5 minutes)
# Runs every 5 minutes
*/5 * * * * root find /tmp -name 'podman_logs_*.log' -mmin +5 -exec sh -c 'session=$(basename {} .log); kill $(ps | grep "curl.*$session" | grep -v grep | awk "{print \$1}") 2>/dev/null; rm -f /tmp/$session.* 2>/dev/null' \;
